# installers/linux/nhibernaut-dashboard.service
[Unit]
Description=NHibernaut Dashboard
# Order after the network stack is up (enough to bind a listening socket). We do NOT
# Want/After network-online.target: the dashboard only *listens* on a local port, it
# makes no outbound connections at startup, and waiting for network-online stalls the
# start ~90s (often to the wait-online timeout) on minimal/CI systemd, tripping the
# default 90s start timeout.
After=network.target

[Service]
Type=notify
ExecStart=/usr/lib/nhibernaut-dashboard/nhibernaut-dashboard
EnvironmentFile=-/etc/nhibernaut-dashboard/dashboard.env
Restart=on-failure
RestartSec=3
DynamicUser=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
